# package/tlt-logtrigger/src/files/blocklist.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Teltonika

START=11
STOP=98

start() {
	logtrigger &
}

stop() {
	killall logtrigger
}
