# package/base-files/files/etc/uci-defaults/30_crt_openvpn
#!/bin/sh
/sbin/crt_openvpn.sh r &

exit 0
