# package/base-files/files/etc/init.d/sim_conf_switch
#!/bin/sh /etc/rc.common
# Copyright (C) 2014 Teltonika

START=43

start() {
	/usr/sbin/sim_conf_switch.sh default
}
