# package/base-files/files/etc/hotplug.d/iface/99-static_arp.sh
if [ "$ACTION" == "ifup" ]; then
	/etc/init.d/static_arp restart
fi

